# apt-mirror-web.conf
description "apt-mirror-web"
author "Dmitrijs Ledkovs <xnox@ubuntu.com>"

start on runlevel [2345]
stop on runlevel [06]

exec uwsgi --plugin python3 --http-socket :80 --wsgi-file apt-mirror-web
